# .config/systemd/user/hypr-lid.path
[Unit]
Description=Watch lid events

[Path]
PathChanged=/proc/acpi/button/lid/LID0/state

[Install]
WantedBy=default.target
