# .config/systemd/user/hypr-lid.path
[Unit]
Description=Watch lid events

[Path]
PathChanged=/run/systemd/seats/

[Install]
WantedBy=default.target

