# .config/systemd/user/hypr-lid.service
[Unit]
Description=Hyprland lid switch handler
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.config/hypr/scripts/lid-handler.sh